# modbus_rtu/etc/init.d/modbus_rtu.service
[Unit]
Description=modbus_rtu

[Service]
ExecStart=/usr/local/sbin/modbus_rtu
Restart=always

[Install]
WantedBy=multi-user.target